# templates/container-splatnet2statink.timer.j2
[Unit]
Description=Runs Podman container-splatnet2statink-{{ name }}.service

[Timer]
OnBootSec=5min
OnUnitActiveSec={{ delay }}

[Install]
WantedBy=timers.target
